# templates/serverApp.service
[Unit]
Description=Nodejs server running
[Service]
ExecStart=/usr/bin/env node /var/www/server.js
Restart=always
[Install]
WantedBy=multi-user.target